# modian-live-wrapper/roles/legacy/templates/modian-legacy.target
[Unit]
Description=Start Modian Legacy system
BindsTo=graphical-session.target
