# modian-live-wrapper/roles/modian-system-modes/templates/user-unit.target
[Unit]
Description={{name}} running
BindsTo=graphical-session.target
